# backend/deploy/ai-secretary-api.service
# AI Secretary API サーバ（POST /devices, GET /briefings/latest）を常駐させる systemd unit。
# 導入手順は docs/specs/deploy-g3plus.md を参照。
#
# 1. User / WorkingDirectory / ExecStart のパスをこのサーバの環境に合わせて書き換える
#    （npm の絶対パスは `command -v npm` で確認。nvm なら ~/.nvm/versions/node/vXX/bin/npm）
# 2. sudo cp backend/deploy/ai-secretary-api.service /etc/systemd/system/
# 3. sudo systemctl daemon-reload && sudo systemctl enable --now ai-secretary-api
# 4. ログ確認: journalctl -u ai-secretary-api -f

[Unit]
Description=AI Secretary API server (devices / briefings)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=akiraak
WorkingDirectory=/home/akiraak/ai-secretary/backend
# .env は WorkingDirectory（backend/）から dotenv が読み込む
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
